# config/hrm_service.conf
[program:hrm-service]
command=/home/yaser/Ultimate skill/HRM/.venv/bin/python -m orchestrator.hrm_service --port 9501
directory=/home/yaser/Ultimate skill/raphael-2.0
user=yaser
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/hrm_service.err.log
stdout_logfile=/var/log/hrm_service.out.log
environment=PYTHONUNBUFFERED=1
